# ansible/roles/watchdogs/files/node-exporter-watchdog.service
[Unit]
Description=Node Exporter health check watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/node-exporter-watchdog.sh
