# ansible/roles/watchdogs/files/prometheus-watchdog.service
[Unit]
Description=Prometheus health check watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/prometheus-watchdog.sh
